# systemd/network-traffic-dashboard.service
[Unit]
Description=Network Traffic Analyzer - web dashboard
After=network.target

[Service]
Type=simple
User=netdashboard
Group=nettraffic
WorkingDirectory=/opt/network-traffic-analyzer
EnvironmentFile=/etc/network-traffic-analyzer/dashboard.env
ExecStart=/opt/network-traffic-analyzer/venv/bin/python3 /opt/network-traffic-analyzer/webapp.py
Restart=on-failure
RestartSec=5

# No packet-capture privileges needed at all: just reads report files and
# serves HTTP, same reasoning as the dashboard container in Docker.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/network-traffic-analyzer/state
# Reports themselves are never written by the dashboard - read access to
# them (and to the state dir above) is via shared nettraffic group
# membership, not an explicit ReadWritePaths grant.

# Further lockdown: a plain Flask app reading report files and serving
# HTTP has no legitimate reason to touch any of the kernel surfaces or
# namespace/privilege-escalation primitives below, so all of them are
# denied outright rather than left available-but-unused.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
# Only what a Flask app doing plain HTTP + local file I/O ever needs -
# AF_UNIX because glibc's resolver can use it to reach a local
# nss/systemd-resolved socket for DNS.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
